# supervisord.conf
; supervisord.conf — process manager for single-container e-flow deployments.
;
; When running `docker run e-flow` (or any standalone deployment), supervisord
; starts and supervises BOTH the data-collection monitor and the Streamlit
; dashboard so operators don't need to manage two separate processes.
;
; In a docker-compose multi-container setup the individual services override
; the CMD to `python monitor.py` / `streamlit run app.py` directly, so this
; file is NOT used in that scenario.

[supervisord]
nodaemon=true
logfile=/app/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid

; ─── Data-collection monitor ───────────────────────────────────────────────
[program:monitor]
command=python monitor.py
directory=/app
autostart=true
autorestart=true
startretries=5
; Write both stdout and stderr to monitor.log (same file monitor.py uses for
; its own RotatingFileHandler so all output is in one place).
stdout_logfile=/app/monitor.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/app/monitor.log
stderr_logfile_maxbytes=0
redirect_stderr=true
environment=MONITOR_ENABLED="true",MONITOR_INTERVAL="60",SCRAPER_FORCE_REQUESTS="1",STORE_ALL_READINGS="false",EXIT_ON_UNHEALTHY="true"

; ─── Streamlit dashboard ────────────────────────────────────────────────────
[program:streamlit]
command=streamlit run app.py --server.port=8501 --server.address=0.0.0.0
directory=/app
autostart=true
autorestart=true
startretries=3
; Stream Streamlit's stdout/stderr straight to the container stdout/stderr so
; `docker logs` shows live output.
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=STREAMLIT_SERVER_PORT="8501",STREAMLIT_SERVER_ADDRESS="0.0.0.0"
